# systemd/ndimon-finder.service
[Unit]
Description=NDIMon-R Source Finder
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
WatchdogSec=60
Restart=always
RestartSec=5
ExecStart=/usr/local/bin/ndimon-finder
Environment=LD_LIBRARY_PATH=/usr/local/lib

[Install]
WantedBy=multi-user.target
